# ethanol-mix.service
[Unit]
Description=Ethanol Mix Calculator
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node __APP_DIR__/server.js
WorkingDirectory=__APP_DIR__
Restart=on-failure
RestartSec=5
Environment=PORT=3000
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
